# root/etc/sv/fan-mode/finish
#!/bin/sh

cacheFile=/var/cache/fan-mode
fanFile=/sys/devices/platform/asus-nb-wmi/fan_boost_mode

# Cache current fan-mode value
[ -r $fanFile ] && echo $(cat $fanFile) > $cacheFile
