# root/etc/sv/fan-mode/run
#!/bin/sh

asus_fanmode 2>&1 > /dev/null
