# root/etc/sv/asus_battery/run
#!/bin/sh

# Set a charging threshold for battery on ASUS laptops.
echo 60 > /sys/class/power_supply/BAT0/charge_control_end_threshold
exec chpst -b asus_battery pause
