# root/etc/sv/nvidia-pm/run
#!/bin/sh

# Remove NVIDIA GPU audio device for runtime D3 power management
GPU_AUDIO=0000:01:00.1
echo 1 > /sys/bus/pci/devices/$GPU_AUDIO/remove
exec chpst -b nvidia-pm pause
